# config/launchd/com.voiceinput.stt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voiceinput.stt</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shifengzhang/anaconda3/envs/vif-stt/bin/python</string>
        <string>/Users/shifengzhang/voice-input-framework/services/stt_server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/shifengzhang/voice-input-framework</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>VIF_LOG_LEVEL</key>
        <string>INFO</string>
        <key>VIF_STT_HOST</key>
        <string>127.0.0.1</string>
        <key>VIF_STT_PORT</key>
        <string>6544</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/shifengzhang/voice-input-framework/logs/stt-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shifengzhang/voice-input-framework/logs/stt-stderr.log</string>
</dict>
</plist>
